# pkg/storeops-basic-custom-application/DEBIAN/prerm
#!/bin/bash

systemctl daemon-reload
systemctl disable storeops-basic-custom-application
systemctl stop storeops-basic-custom-application
