# pkg/storeops-custom-application/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable storeops-custom-application
systemctl restart storeops-custom-application
rm -rf /usr/share/cockpit/storeops-custom-application
