# pkg/storeops-custom-application/DEBIAN/prerm
#!/bin/bash

systemctl daemon-reload
systemctl disable storeops-custom-application
systemctl stop storeops-custom-application
