# pkg/APPLICATION_ID/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable <APPLICATION_ID>
systemctl restart <APPLICATION_ID>
