# pkg/APPLICATION_ID/DEBIAN/prerm
#!/bin/bash

systemctl daemon-reload
systemctl disable <APPLICATION_ID>
systemctl stop <APPLICATION_ID>
